# packaging/scripts/postinstall
#!/bin/sh

systemctl daemon-reload
systemctl enable sysload-exporter.service
systemctl restart sysload-exporter.service