# packaging/scripts/preremove
#!/bin/sh

systemctl stop sysload-exporter.service
systemctl disable sysload-exporter.service